# systemd/mugi@.service
[Unit]
Description=mugi sync
After=network.target
StartLimitIntervalSec=90
StartLimitBurst=3
 
[Service]
User=mirrors
Group=mirrors
ExecStart=mugi /etc/mugi/mugi.conf %i
Restart=on-failure
 
[Install]
WantedBy=multi-user.target
