# systemd/mugi@archlinuxcn.timer
[Unit]
Description=mugi sync archlinuxcn every 5h

[Timer]
OnActiveSec=300s
OnUnitActiveSec=5h

[Install]
WantedBy=multi-user.target
